# attacks/persistence/T1547.001-registry-run-keys/logs/001-reg-add-run-key.json
{
  "_metadata": {
    "attack_id": "T1547.001",
    "variation": "001",
    "name": "Registry Run Keys Persistence",
    "description": "Add registry run key for persistence via HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
    "expected_detection": "HIGH",
    "kill_chain_phase": "persistence",
    "format": "wazuh-native",
    "format_version": "1.0",
    "campaign": "Campaign 21",
    "date_executed": "2026-02-09"
  },
  "logs": [
    {
      "_comment": "Sysmon Event 1 - Process Creation: reg add Run key",
      "sequence": 1,
      "timestamp": "2026-02-09T10:42:00.000Z",
      "log": {
        "timestamp": "",
        "rule": {
          "level": 12,
          "description": "Sysmon - Event 1: Process creation",
          "id": "200151",
          "mitre": {
            "id": ["T1547.001"],
            "tactic": ["Persistence", "Privilege Escalation"],
            "technique": ["Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder"]
          },
          "firedtimes": 1,
          "mail": false,
          "groups": ["sysmon", "sysmon_event1", "windows"]
        },
        "agent": {"id": "014", "name": "WS01", "ip": "10.1.0.20"},
        "manager": {"name": "wazuh.manager"},
        "id": "1709000106.000001",
        "decoder": {"name": "windows_eventchannel"},
        "location": "EventChannel",
        "data": {
          "win": {
            "system": {
              "eventID": "1",
              "channel": "Microsoft-Windows-Sysmon/Operational",
              "providerName": "Microsoft-Windows-Sysmon",
              "computer": "WS01.yourcompany.local",
              "level": "4",
              "task": "1",
              "severityValue": "INFORMATION"
            },
            "eventdata": {
              "utcTime": "2026-02-09 10:42:00.000",
              "processGuid": "{12345678-1234-1234-1234-000000106001}",
              "processId": "9256",
              "image": "C:\\Windows\\System32\\reg.exe",
              "commandLine": "reg  add \"HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\" /v Campaign21 /t REG_SZ /d \"C:\\Windows\\System32\\cmd.exe /c echo persistence\" /f",
              "currentDirectory": "C:\\Windows\\System32\\",
              "user": "YOURCOMPANY\\Administrator",
              "integrityLevel": "High",
              "hashes": "SHA256=B8C9D0E1F2A3B4C5D6E7F8A9B0C1D2E3F4A5B6C7D8E9F0A1B2C3D4E5F6A7B8C9",
              "parentProcessGuid": "{12345678-1234-1234-1234-000000106000}",
              "parentProcessId": "1234",
              "parentImage": "C:\\Windows\\System32\\cmd.exe",
              "parentCommandLine": "C:\\Windows\\System32\\cmd.exe"
            }
          }
        }
      }
    },
    {
      "_comment": "Sysmon Event 13 - Registry Value Set",
      "sequence": 2,
      "timestamp": "2026-02-09T10:42:01.000Z",
      "log": {
        "timestamp": "",
        "rule": {
          "level": 12,
          "description": "Sysmon - Event 13: Registry value set",
          "id": "200157",
          "mitre": {
            "id": ["T1547.001"],
            "tactic": ["Persistence", "Privilege Escalation"],
            "technique": ["Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder"]
          },
          "firedtimes": 1,
          "mail": false,
          "groups": ["sysmon", "sysmon_event13", "windows"]
        },
        "agent": {"id": "014", "name": "WS01", "ip": "10.1.0.20"},
        "manager": {"name": "wazuh.manager"},
        "id": "1709000106.000002",
        "decoder": {"name": "windows_eventchannel"},
        "location": "EventChannel",
        "data": {
          "win": {
            "system": {
              "eventID": "13",
              "channel": "Microsoft-Windows-Sysmon/Operational",
              "providerName": "Microsoft-Windows-Sysmon",
              "computer": "WS01.yourcompany.local",
              "level": "4",
              "task": "13",
              "severityValue": "INFORMATION"
            },
            "eventdata": {
              "utcTime": "2026-02-09 10:42:01.000",
              "processGuid": "{12345678-1234-1234-1234-000000106001}",
              "processId": "9256",
              "image": "C:\\Windows\\System32\\reg.exe",
              "targetObject": "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Campaign21",
              "details": "C:\\Windows\\System32\\cmd.exe /c echo persistence",
              "eventType": "SetValue"
            }
          }
        }
      }
    }
  ]
}
