# attacks/persistence/T1053.005-scheduled-task/logs/001-schtasks-create-system.json
{
  "_metadata": {
    "attack_id": "T1053.005",
    "variation": "001",
    "name": "Scheduled Task Creation via schtasks for SYSTEM",
    "description": "Attacker creates scheduled task running as SYSTEM for persistence",
    "expected_detection": "HIGH",
    "kill_chain_phase": "persistence",
    "format": "wazuh-native",
    "format_version": "1.0",
    "campaign": "operation-shadow-domain"
  },
  "logs": [
    {
      "_comment": "Sysmon Event 1 - schtasks /create",
      "sequence": 1,
      "timestamp": "2026-02-04T19:48:00.000Z",
      "log": {
        "timestamp": "2026-02-04T19:48:00.000Z",
        "rule": {
          "level": 3,
          "description": "Sysmon - Event 1: Process Create",
          "id": "92000",
          "firedtimes": 1,
          "mail": false,
          "groups": ["sysmon", "sysmon_event1"]
        },
        "agent": {"id": "003", "name": "DC01", "ip": "10.0.10.10"},
        "manager": {"name": "wazuh.manager"},
        "id": "1706900002.000001",
        "decoder": {"name": "windows_eventchannel"},
        "location": "EventChannel",
        "data": {
          "win": {
            "system": {
              "eventID": "1",
              "channel": "Microsoft-Windows-Sysmon/Operational",
              "providerName": "Microsoft-Windows-Sysmon",
              "computer": "DC01.Pooka.local"
            },
            "eventdata": {
              "utcTime": "2026-02-04 19:48:00.000",
              "processGuid": "{DC01-0002-0001-0001-000000000001}",
              "processId": "6789",
              "image": "C:\\Windows\\System32\\schtasks.exe",
              "commandLine": "schtasks /create /tn ShadowTask /tr cmd.exe /sc ONLOGON /ru SYSTEM",
              "currentDirectory": "C:\\Users\\attacker\\",
              "user": "POOKA\\attacker",
              "integrityLevel": "High",
              "parentProcessGuid": "{DC01-0002-0001-0001-000000000000}",
              "parentProcessId": "1234",
              "parentImage": "C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe",
              "parentCommandLine": "powershell.exe"
            }
          }
        }
      }
    },
    {
      "_comment": "Security Event 4698 - Scheduled Task Created",
      "sequence": 2,
      "timestamp": "2026-02-04T19:48:01.000Z",
      "log": {
        "timestamp": "2026-02-04T19:48:01.000Z",
        "rule": {
          "level": 8,
          "description": "Scheduled Task Created",
          "id": "60144",
          "mitre": {
            "id": ["T1053.005"],
            "tactic": ["Persistence", "Privilege Escalation"],
            "technique": ["Scheduled Task/Job: Scheduled Task"]
          },
          "firedtimes": 1,
          "mail": false,
          "groups": ["windows", "scheduled_task"]
        },
        "agent": {"id": "003", "name": "DC01", "ip": "10.0.10.10"},
        "manager": {"name": "wazuh.manager"},
        "id": "1706900002.000002",
        "decoder": {"name": "windows_eventchannel"},
        "location": "EventChannel",
        "data": {
          "win": {
            "system": {
              "eventID": "4698",
              "channel": "Security",
              "providerName": "Microsoft-Windows-Security-Auditing",
              "computer": "DC01.Pooka.local"
            },
            "eventdata": {
              "subjectUserName": "attacker",
              "subjectDomainName": "POOKA",
              "taskName": "\\ShadowTask",
              "taskContent": "cmd.exe"
            }
          }
        }
      }
    }
  ]
}
