# attacks/persistence/T1053.005-scheduled-task/logs/001-schtasks-persistence.json
{
  "_metadata": {
    "attack_id": "T1053.005",
    "variation": "001",
    "name": "Scheduled Task Persistence",
    "description": "Create persistent scheduled task using schtasks.exe",
    "expected_detection": "MEDIUM",
    "kill_chain_phase": "persistence",
    "format": "wazuh-native",
    "format_version": "1.0",
    "campaign": "Campaign 26",
    "date_executed": "2026-02-09"
  },
  "logs": [
    {
      "_comment": "Sysmon Event 1 - Process Creation: schtasks /create",
      "sequence": 1,
      "timestamp": "2026-02-09T15:10:00.000Z",
      "log": {
        "timestamp": "",
        "rule": {
          "level": 10,
          "description": "Sysmon - Event 1: Process creation",
          "id": "200151",
          "mitre": {
            "id": ["T1053.005"],
            "tactic": ["Execution", "Persistence", "Privilege Escalation"],
            "technique": ["Scheduled Task/Job: Scheduled Task"]
          },
          "firedtimes": 1,
          "mail": true,
          "groups": ["sysmon", "sysmon_event1", "windows"]
        },
        "agent": {"id": "012", "name": "DC01", "ip": "10.1.0.10"},
        "manager": {"name": "wazuh.manager"},
        "id": "1709000305.000001",
        "decoder": {"name": "windows_eventchannel"},
        "location": "EventChannel",
        "data": {
          "win": {
            "system": {
              "eventID": "1",
              "channel": "Microsoft-Windows-Sysmon/Operational",
              "providerName": "Microsoft-Windows-Sysmon",
              "computer": "DC01.yourcompany.local",
              "level": "4",
              "task": "1",
              "severityValue": "INFORMATION"
            },
            "eventdata": {
              "utcTime": "2026-02-09 15:10:00.000",
              "processGuid": "{12345678-1234-1234-1234-000000305001}",
              "processId": "29816",
              "image": "C:\\Windows\\System32\\schtasks.exe",
              "commandLine": "schtasks  /create /tn \"Campaign26Task\" /tr \"C:\\Windows\\System32\\cmd.exe /c echo task\" /sc daily /st 00:00",
              "currentDirectory": "C:\\Windows\\System32\\",
              "user": "YOURCOMPANY\\Administrator",
              "integrityLevel": "High",
              "hashes": "SHA256=7A8B9C0D1E2F3A4B5C6D7E8F9A0B1C2D3E4F5A6B7C8D9E0F1A2B3C4D5E6F7A8B",
              "parentProcessGuid": "{12345678-1234-1234-1234-000000305000}",
              "parentProcessId": "1234",
              "parentImage": "C:\\Windows\\System32\\cmd.exe",
              "parentCommandLine": "C:\\Windows\\System32\\cmd.exe"
            }
          }
        }
      }
    },
    {
      "_comment": "Windows Event 4698 - Scheduled Task Created",
      "sequence": 2,
      "timestamp": "2026-02-09T15:10:01.000Z",
      "log": {
        "timestamp": "",
        "rule": {
          "level": 10,
          "description": "Scheduled task created",
          "id": "60116",
          "mitre": {
            "id": ["T1053.005"],
            "tactic": ["Execution", "Persistence", "Privilege Escalation"],
            "technique": ["Scheduled Task/Job: Scheduled Task"]
          },
          "firedtimes": 1,
          "mail": true,
          "groups": ["windows", "windows_security"]
        },
        "agent": {"id": "012", "name": "DC01", "ip": "10.1.0.10"},
        "manager": {"name": "wazuh.manager"},
        "id": "1709000305.000002",
        "decoder": {"name": "windows_eventchannel"},
        "location": "EventChannel",
        "data": {
          "win": {
            "system": {
              "eventID": "4698",
              "channel": "Security",
              "providerName": "Microsoft-Windows-Security-Auditing",
              "computer": "DC01.yourcompany.local",
              "level": "0",
              "task": "12804",
              "severityValue": "AUDIT_SUCCESS"
            },
            "eventdata": {
              "SubjectUserSid": "S-1-5-21-123456789-123456789-123456789-500",
              "SubjectUserName": "Administrator",
              "SubjectDomainName": "YOURCOMPANY",
              "TaskName": "\\Campaign26Task",
              "TaskContent": "<?xml version=\"1.0\" encoding=\"UTF-16\"?><Task>...</Task>"
            }
          }
        }
      }
    }
  ]
}
